# mealie-sync.service
[Unit]
Description=Mealie and Google Calendar Two-Way Sync
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/mealie_integration
ExecStart=/usr/bin/node /opt/mealie_integration/sync_master.js
# This ensures logs are sent to journalctl immediately
Environment=NODE_OPTIONS=--no-warnings
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
